# services/com.fountainbuddy.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fountainbuddy.service</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/adam/Documents/Projects/fountain-buddy/venv/bin/python</string>
        <string>/Users/adam/Documents/Projects/fountain-buddy/run.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/adam/Documents/Projects/fountain-buddy</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/adam/Documents/Projects/fountain-buddy/logs/fountain-buddy.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/adam/Documents/Projects/fountain-buddy/logs/fountain-buddy-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/adam/Documents/Projects/fountain-buddy</string>
    </dict>
    
    <key>StartInterval</key>
    <integer>30</integer>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>