# services/com.fountainbuddy.photocleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fountainbuddy.photocleanup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/adam/Documents/Projects/fountain-buddy/venv/bin/python</string>
        <string>/Users/adam/Documents/Projects/fountain-buddy/photo_cleanup_scheduler.py</string>
        <string>daily</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/adam/Documents/Projects/fountain-buddy</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/adam/Documents/Projects/fountain-buddy/logs/photocleanup.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/adam/Documents/Projects/fountain-buddy/logs/photocleanup.error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>